# deploy/systemd/feature-service.service
[Unit]
Description=AIValid Feature Service (.NET 8)
After=network.target

[Service]
Type=simple
User=aivalid
Group=aivalid
WorkingDirectory=/opt/aivalid/feature-service
EnvironmentFile=/opt/aivalid/feature-service/.env
ExecStart=/usr/bin/dotnet /opt/aivalid/feature-service/FeatureService.Api.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=feature-service
LimitNOFILE=65535

# --- systemd hardening (defense-in-depth) ---
# NOTE: Do NOT enable MemoryDenyWriteExecute=true here because .NET uses a JIT.
NoNewPrivileges=true
PrivateTmp=true
PrivateDevices=true
ProtectSystem=full
ProtectHome=true
ProtectControlGroups=true
ProtectKernelModules=true
ProtectKernelTunables=true
LockPersonality=true
RestrictSUIDSGID=true
RestrictRealtime=true
SystemCallArchitectures=native
UMask=0077
CapabilityBoundingSet=
AmbientCapabilities=

[Install]
WantedBy=multi-user.target
